# deploy/restaurant_orders.socket
[Unit]
Description=Gunicorn socket for restaurant_orders

[Socket]
ListenStream=/run/restaurant_orders.sock

[Install]
WantedBy=sockets.target